# deploy/systemd/helix-gateway@.service
[Unit]
Description=Helix API Gateway (profile %i)
Documentation=https://github.com/javded-itres/HelixAgent/blob/master/docs/en/DEPLOYMENT.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helix
Group=helix
WorkingDirectory=/home/helix
Environment=HOME=/home/helix
Environment=HELIX_HOME=/home/helix/.helix
EnvironmentFile=-/etc/helix/helix.conf
# Secrets and gateway bind: ~/.helix/profiles/%i/.env (loaded by gateway_worker)
ExecStart=${HELIX_PYTHON} -m cli.services.gateway_worker --profile %i
ExecStop=${HELIX_BIN} -p %i gateway stop
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target